# initscript/stap@.service.in
[Unit]
Description=SystemTap script %i
Documentation=man:stap(1)
After=local-fs.target
Before=staprun@%i.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/systemtap/script.d/%i.options
ExecStartPre=@prefix@/libexec/systemtap/stap-service-prepare %i
ExecStart=/bin/systemctl --no-block start staprun@%i.service
ExecStop=/bin/systemctl stop staprun@%i.service
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
